# services/activitywatch.service
[Unit]
Description=Activitywatch
After=network.target

[Service]
Type=simple
ExecStart=/home/kage/bin/activitywatch/aw-server

