# services/auto-shutdown.service
[Unit]
Description=Shutdown automatically at 10:30PM everyday

[Service]
Type=simple
ExecStart=/sbin/shutdown -h +2
